# docker/cron/crontab
# cronjob
MAILTO = ""
SHELL=/bin/sh
BASH_ENV=/etc/profile
33 3 * * * sh -c 'cd /app && /usr/local/bin/python3 /app/data_grabber.py'



